# ops/launchd/com.omnyist.questlog.poll-steam-warframe.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Runs every 5 minutes on Saya's macOS host. Polls Steam for Warframe
  play state; when a session ends, inline-triggers archive_warframe.

  Install:
    cp ops/launchd/com.omnyist.questlog.poll-steam-warframe.plist \
       ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.omnyist.questlog.poll-steam-warframe.plist

  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.omnyist.questlog.poll-steam-warframe.plist
    rm ~/Library/LaunchAgents/com.omnyist.questlog.poll-steam-warframe.plist

  Logs: /tmp/questlog-poll-steam-warframe.{out,err}.log
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.omnyist.questlog.poll-steam-warframe</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/docker</string>
        <string>exec</string>
        <string>questlog-server-1</string>
        <string>uv</string>
        <string>run</string>
        <string>python</string>
        <string>manage.py</string>
        <string>poll_steam_warframe</string>
    </array>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/questlog-poll-steam-warframe.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/questlog-poll-steam-warframe.err.log</string>
</dict>
</plist>
